# agx-startup.service
[Unit]
Description=Mika's startup script

[Service]
ExecStart=/usr/local/bin/agx-startup.py

[Install]
WantedBy=multi-user.target
