# services/iperf.service
[Unit]
Description=iperf server

[Service]
Type=simple
Restart=always
ExecStart=/usr/bin/iperf -s

[Install]
WantedBy=multi-user.target
