# wacom/99-wacom-tablet-android.rules
# For the Android mode of the tablet
SUBSYSTEM=="usb", \
ATTRS{idVendor}=="2d1f", \
ATTRS{idProduct}=="0378", \
RUN+="/usr/local/bin/wacom.sh"
